# contrib/dhcp-adguard-sync.service
[Unit]
Description=DHCP to AdGuard Home Sync Service
After=network.target

[Service]
ExecStart=/usr/bin/dhcp-adguard-sync serve
Restart=always
User=root
EnvironmentFile=/etc/dhcp-adguard-sync/config.yaml

[Install]
WantedBy=multi-user.target